# power_button.service
[Unit]
Description=power button
After=syslog.service

[Service]
ExecStart=/opt/power_button/power_button.sh

[Install]
WantedBy=multi-user.target
